# .config/systemd/user/agentctl-project-index.timer
[Unit]
Description=agentctl schedule: project-index
# Generated by 'agentctl reload' — do not edit manually.

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
Unit=agentctl@project-index.service

[Install]
WantedBy=timers.target
