# .config/systemd/user/agentctl-lab-sync.timer
[Unit]
Description=agentctl schedule: lab-sync
# Generated by 'agentctl reload' — do not edit manually.

[Timer]
OnCalendar=Mon *-*-* 04:00:00
Persistent=true
Unit=agentctl@lab-sync.service

[Install]
WantedBy=timers.target
