# .config/systemd/user/agentctl-nightly-sync.timer
[Unit]
Description=agentctl schedule: nightly-sync
# Generated by 'agentctl reload' — do not edit manually.

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true
Unit=agentctl@nightly-sync.service

[Install]
WantedBy=timers.target
