# .config/systemd/user/fleet-pulse.timer
[Unit]
Description=Push fleet-pulse heartbeat every minute

[Timer]
OnCalendar=*:0/1
AccuracySec=10
Persistent=true

[Install]
WantedBy=timers.target
