# .config/systemd/user/git-sync-notes.timer
[Unit]
Description=Run Git Auto Push and Pull periodically

[Timer]
OnCalendar=hourly
AccuracySec=1

[Install]
WantedBy=timers.target
