# .config/systemd/user/git-sync-notes.service
[Unit]
Description=Sync ~/.notes through Forgejo (HTTPS) primary and GitHub backup remote
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/.notes
EnvironmentFile=-%h/.config/notes-sync.env
ExecStart=%h/.local/src/git/git-sync-notes.sh

[Install]
WantedBy=default.target
