# .config/systemd/user/git-sync-notes.timer
[Unit]
Description=Run notes sync automatically (5-min fallback for the MQTT/ntfy webhook fan-out)

[Timer]
OnBootSec=2m
# 5-minute fallback. The Forgejo webhook -> notes-sync-bridge -> mosquitto/ntfy
# path is the live trigger; this timer is just a safety net for cases where
# the bridge / broker is down or the device just woke up.
OnUnitActiveSec=5min
RandomizedDelaySec=30s
AccuracySec=15s
Persistent=true
Unit=git-sync-notes.service

[Install]
WantedBy=timers.target
