# .config/systemd/user/inbox-notify.service
[Unit]
Description=Send notification for inbox activity

[Service]
Type=oneshot
ExecStart=%h/.local/bin/inbox-notify
