# .config/systemd/user/infra-collect.service
[Unit]
Description=Infrastructure Dashboard Collector
After=network.target

[Service]
Type=oneshot
ExecStart=%h/.local/src/infra-dash/collect.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
