# .config/systemd/user/journal-refs-archive.service
[Unit]
Description=Archive old journal ref folders monthly

[Service]
Type=oneshot
ExecStart=%h/.local/bin/archive-refs
