# .config/systemd/user/journal-monthly-archive.service
[Unit]
Description=Archive daily journal entries monthly (notes CLI)

[Service]
Type=oneshot
# Roll the previous calendar month into the monthly summary + move dailies to archive.
#   journalctl --user -u journal-monthly-archive
ExecStart=%h/.local/bin/notes --verbose archive
