# .config/systemd/user/journal-daily-summarize.timer
[Unit]
Description=Run journal daily summarization nightly

[Timer]
OnCalendar=*-*-* 01:00:00
Persistent=true

[Install]
WantedBy=timers.target
