# .config/systemd/user/journal-daily-summarize.service
[Unit]
Description=Summarize yesterday's daily journal entry

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/kblack0610/.notes/journal/scripts/daily_summarize_cleanup.py
StandardOutput=append:/home/kblack0610/.notes/journal/scripts/daily_summarize.log
StandardError=append:/home/kblack0610/.notes/journal/scripts/daily_summarize.log
