# .config/systemd/user/notes-mqtt.service
[Unit]
Description=Subscribe to mosquitto notes/sync/needed and trigger notes-sync
After=network-online.target
Wants=network-online.target
ConditionFileIsExecutable=/usr/bin/mosquitto_sub

[Service]
Type=simple
EnvironmentFile=-%h/.config/notes-sync.env
Environment=NOTES_MQTT_HOST=mosquitto.kblab.me
Environment=NOTES_MQTT_PORT=31883
Environment=NOTES_MQTT_TOPIC=notes/sync/needed
ExecStart=/bin/sh -c 'exec mosquitto_sub -h "$NOTES_MQTT_HOST" -p "$NOTES_MQTT_PORT" -t "$NOTES_MQTT_TOPIC" -q 1 --keepalive 60 -v | while read -r _line; do "$HOME/.local/bin/notes-sync" >>"$HOME/.local/state/notes-sync/mqtt.log" 2>&1 || true; done'
Restart=always
RestartSec=10s
Nice=10

[Install]
WantedBy=default.target
