# .config/systemd/user/regen-project-index.path
[Unit]
Description=Watch lab project folders for stage moves (regenerate the project index)
# Watches the STAGE dirs, not index.md — so the regen writing index.md never self-triggers.

[Path]
PathChanged=%h/.notes/lab/projects/current
PathChanged=%h/.notes/lab/projects/prod
PathChanged=%h/.notes/lab/projects/archived
Unit=agentctl@project-index.service

[Install]
WantedBy=default.target
