# .config/systemd/user/theme-day.service
[Unit]
Description=Apply day theme (jackie-brown)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/theme-switch jackie-brown
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1
