# .config/systemd/user/waybar.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.config/waybar/launch.sh
