# .config/systemd/user/sunshine.service.d/override.conf
[Service]
Environment="GBM_BACKEND=nvidia-drm"
Environment="__GLX_VENDOR_LIBRARY_NAME=nvidia"
Environment="WLR_NO_HARDWARE_CURSORS=1"
