# .config/systemd/user/agentctl@sentinel.service.d/override.conf
[Service]
# ~/.local/bin holds agent-notify, aws-cred-timer, watch-companion-loop.
# systemd's default user PATH omits it, so probes/notify silently fail without this.
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
