# .config/launchd/com.kblack.git-sync-notes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kblack.git-sync-notes</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>~/.local/bin/notes-sync</string>
    </array>

    <!-- 5-minute fallback. The MQTT/ntfy webhook fan-out catches the live
         path; this timer just makes sure no edit gets stuck if the bridge
         or broker is down for a bit. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/git-sync-notes.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/git-sync-notes.err.log</string>
</dict>
</plist>
